# packaging/systemd/uconsole-webdash.service
[Unit]
Description=uConsole Web Dashboard
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/uconsole/webdash/webdash.py
Environment=WEBDASH_PORT=8080
WorkingDirectory=/opt/uconsole/webdash
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
