# packaging/systemd/uconsole-status.service
[Unit]
Description=Push uConsole status to uconsole.cloud
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/uconsole/scripts/system/push-status.sh
EnvironmentFile=/etc/uconsole/status.env
User=UCONSOLE_USER
Group=UCONSOLE_USER
