# packaging/systemd/uconsole-update.timer
[Unit]
Description=uConsole weekly system update (Sunday 4am)

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
