# packaging/systemd/uconsole-webdash.service
[Unit]
Description=uConsole Web Dashboard
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/uconsole/webdash/app.py
Environment=WEBDASH_PORT=8080
WorkingDirectory=/opt/uconsole/webdash
User=UCONSOLE_USER
Group=UCONSOLE_USER
Restart=always
RestartSec=3
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
