# device/scripts/config/systemd-user/claude-relay-agent.service
[Unit]
Description=Claude Relay Agent — polls cloud broker for Telegram messages
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/claude-relay-repo/relay-agent.sh
Restart=on-failure
RestartSec=10
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
